# docs/operations/systemd.service.example
[Unit]
Description=Sand Rocket API Server
After=network.target

[Service]
Type=simple
User=sandrocket
WorkingDirectory=/home/sandrocket/sandrocket
Environment=NODE_ENV=production
EnvironmentFile=/home/sandrocket/sandrocket/.env
ExecStart=/home/sandrocket/.nvm/versions/node/v20.11.1/bin/node /home/sandrocket/sandrocket/server.js
# Note: server.js now uses tsx internally to run the TypeScript source,
# which handles workspace module resolution correctly
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

